# vendor/github.com/coreos/rkt/stage1/units/units/shutdown.service
[Unit]
Description=Pod shutdown
AllowIsolate=true
StopWhenUnneeded=yes
DefaultDependencies=false

[Service]
RemainAfterExit=yes
ExecStop=/usr/bin/systemctl halt --force
