# deploy/addons/auto-pause/auto-pause.service
[Unit]
Description=Auto Pause Service

[Service]
Type=simple
ExecStart=/usr/local/bin/auto-pause
Restart=always

[Install]
WantedBy=multi-user.target