# deploy/iso/minikube-iso/arch/aarch64/package/buildkit-bin-aarch64/buildkit.socket
[Unit]
Description=BuildKit
Documentation=https://github.com/moby/buildkit

[Socket]
ListenStream=%t/buildkit/buildkitd.sock
SocketMode=0660
SocketUser=root
SocketGroup=buildkit

[Install]
WantedBy=sockets.target
