# hack/jenkins/cron/Darwin.crontab
*/20 * * * * /Users/jenkins/cleanup-and-reboot.sh
