# deploy/kicbase/nerdctld/nerdctld.socket
[Unit]
Description=nerdctld

[Socket]
ListenStream=/var/run/nerdctld.sock
SocketMode=0666

[Install]
WantedBy=sockets.target