# deploy/kicbase/ipchecker.service
Description=Example systemd service.

[Service]
Type=simple
ExecStart=ipchecker
Restart=always

[Install]
WantedBy=multi-user.target
